# bro-sandbox/sandbox.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user  command
* *     * * *   root    /usr/local/bin/sandbox_disk_limit
* *     * * *   root    /usr/local/bin/sandbox_remove_old_containers 3
* *     * * *   root    /usr/local/bin/sandbox_remove_old_users 3
